# pymultimon/pymultimon.service
[Unit]
Description=PyMultimon
After=network.target

[Service]
ExecStart=/usr/local/bin/pymultimon --syslog

[Install]
WantedBy=multi-user.target
